# supervisord_hive.conf
[supervisord]
nodaemon=true
user=root

[program:hive]
command=/usr/bin/bash -c "/root/run_hive.sh"
stopsignal=TERM